# ops/systemd/tls-pq-bench-track.timer
[Unit]
Description=Run TLS PQ recurring tracking benchmark every 4 hours

[Timer]
OnCalendar=*-*-* 00,04,08,12,16,20:00:00
Persistent=true
Unit=tls-pq-bench-track.service

[Install]
WantedBy=timers.target
